# DEBIAN/postinst
#!/bin/bash

LONG_DATE=`date +%Y%m%d%H%M%S`

chmod +x /etc/init.d/d1-index-task-generator
chmod +x /etc/init.d/d1-index-task-processor

update-rc.d d1-index-task-generator defaults 99
update-rc.d d1-index-task-processor defaults 99

#provide context properties files for postgres, log4j, solr location
cp /usr/share/dataone-cn-index/debian/log4j-generator.properties /etc/dataone/index/log4j-generator.properties
cp /usr/share/dataone-cn-index/debian/log4j-processor.properties /etc/dataone/index/log4j-processor.properties
cp /usr/share/dataone-cn-index/debian/jdbc.properties /etc/dataone/index/jdbc.properties
cp /usr/share/dataone-cn-index/debian/solr.properties /etc/dataone/index/solr.properties

#provide metadata index generation configuration files
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-eml-base.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-eml200.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-eml201.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-eml210.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-eml211.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-resourcemap.xml /etc/dataone/index/index-generation-context
cp /usr/share/dataone-cn-index/debian/index-generation-context/application-context-systemmeta100.xml /etc/dataone/index/index-generation-context

/usr/share/dataone-cn-index/debian/index-generation-context/application-context-fgdc-base.xml /etc/dataone/index/index-generation-context
/usr/share/dataone-cn-index/debian/index-generation-context/application-context-fgdc-std-001-1998.xml /etc/dataone/index/index-generation-context
/usr/share/dataone-cn-index/debian/index-generation-context/application-context-fgdc-std-0011-1999.xml /etc/dataone/index/index-generation-context
/usr/share/dataone-cn-index/debian/index-generation-context/application-context-fgdc-std-0012-1999.xml /etc/dataone/index/index-generation-context

chown -R tomcat6.tomcat6 /etc/dataone/index

cp /usr/share/dataone-cn-index/d1_index_task_generator_daemon.jar /usr/share/java
cp /usr/share/dataone-cn-index/d1_index_task_processor_daemon.jar /usr/share/java
cp /usr/share/dataone-cn-index/d1_index_build_tool.jar /usr/share/java

chown -R tomcat6.tomcat6 /usr/share/java/d1_index_task_generator_daemon.jar
chown -R tomcat6.tomcat6 /usr/share/java/d1_index_task_processor_daemon.jar
chown -R tomcat6.tomcat6 /usr/share/java/d1_index_build_tool.jar

mkdir  /var/log/dataone/index
chown tomcat6.tomcat6 /var/log/dataone/index


###############################################################################
# Configure Postgres
###############################################################################
PG=postgresql-8.4
PG_CONF=/etc/postgresql/8.4/main
PG_USER=postgres

INDEX_DB=d1-index-queue
INDEX_DB_USER=indexer
INDEX_DB_PASS=indexer

## modify pg_hba.conf
PG_HBA_IS_MODIFIED=`grep "${INDEX_DB} ${INDEX_DB_USER}" ${PG_CONF}/pg_hba.conf`
if [ "${PG_HBA_IS_MODIFIED}" == "" ]
then
  echo "backing up ${PG_CONF}/pg_hba.conf to ${PG_CONF}/pg_hba.conf.bak.index"
  cp ${PG_CONF}/pg_hba.conf ${PG_CONF}/pg_hba.conf.bak.index
  chown ${PG_USER} ${PG_CONF}/pg_hba.conf
  chgrp ${PG_USER} ${PG_CONF}/pg_hba.conf

  echo "appending 'host ${INDEX_DB} ${INDEX_DB_USER} 127.0.0.1 255.255.255.255 password' to ${PG_CONF}/pg_hba.conf"
  echo "host ${INDEX_DB} ${INDEX_DB_USER} 127.0.0.1 255.255.255.255 password" >> ${PG_CONF}/pg_hba.conf
fi

## create index task schema and user
echo "Creating index task database schema"
su ${PG_USER} -c "createdb ${INDEX_DB}"

echo "Creating index task db user - indexer" 
su ${PG_USER} -c "psql -c \"CREATE USER ${INDEX_DB_USER} WITH UNENCRYPTED PASSWORD '${INDEX_DB_PASS}'\""


## Restart the postgres db
echo "Restarting postgres database"
/etc/init.d/${PG} restart
######################################################

#open up the correct port for hazelcast
ufw allow 5702


######## Update solr index schema
/etc/init.d/tomcat6 stop

mkdir  /etc/dataone/index/solr

mv /usr/share/solr/d1-cn-index/conf/schema.xml /usr/share/solr/d1-cn-index/conf/schema.xml.bak
cp /usr/share/dataone-cn-index/debian/index-solr-schema.xml /etc/dataone/index/solr/schema.xml

ln -s /etc/dataone/index/solr/schema.xml /usr/share/solr/d1-cn-index/conf/schema.xml

cp /usr/share/dataone-cn-index/debian/d1search.xsl /etc/dataone/index/solr/
ln -s /etc/dataone/index/solr/d1search.xsl /usr/share/solr/d1-cn-index/conf/xslt/d1search.xsl

chown -R tomcat6.tomcat6 /etc/dataone/index/solr
chown -R tomcat6.tomcat6 /usr/share/solr/d1-cn-index/

/etc/init.d/tomcat6 start

######### Finished! inform user how to start daemons
echo ""
echo "Build or re-build Solr Index with /usr/share/dataone-cn-index/d1_index_build_tool.jar while index task processor daemon is not running."
echo "Start Index Task Generator with /etc/init.d/d1-index-task-generator start"
echo "Start Index Task Processor with /etc/init.d/d1-index-task-processor start"

## Update DateONE Version Info Doc
java -jar /usr/share/dataone-cn-version-tool/dataone-cn-version-tool.jar -F/usr/share/dataone-cn-version-tool/version-tool.properties -html > /var/www/cn-version.html
